# deployment/chatwoot3.target
[Unit]
Wants=chatwoot-web.3.service chatwoot-worker.3.service

[Install]
WantedBy=multi-user.targets